# tf1d.service
[Unit]
Description=Traktor F1 userspace driver

[Service]
ExecStart=__INSTALL_PATH__/bin/tf1d

